# cronjob/crontab
30 23 * * * /usr/local/bin/python /app/manage.py import_hotel_data >> /app/cronjob/cron.log 2>&1 # Execute cronjob every day at 23:30
